# scripts/stock-market-app-backup.service
; Nightly pg_dump of the stock-market-app DB.
; Install:
;   sudo cp scripts/stock-market-app-backup.service /etc/systemd/system/
;   sudo cp scripts/stock-market-app-backup.timer   /etc/systemd/system/
;   sudo mkdir -p /var/backups/stock-market-app
;   sudo chown postgres:postgres /var/backups/stock-market-app
;   sudo systemctl daemon-reload
;   sudo systemctl enable --now stock-market-app-backup.timer
;
; Keeps the last 14 dumps; older ones are deleted.

[Unit]
Description=Stock Market App — nightly pg_dump
After=postgresql.service
Requires=postgresql.service

[Service]
Type=oneshot
User=postgres
ExecStart=/bin/bash -c 'umask 077; /usr/bin/pg_dump -Fc stockmarket -f /var/backups/stock-market-app/stockmarket-$(date +%%Y%%m%%d-%%H%%M).dump'
ExecStartPost=/bin/bash -c 'cd /var/backups/stock-market-app && ls -1t stockmarket-*.dump 2>/dev/null | tail -n +15 | xargs -r rm --'
